# launchd/com.youtube-monitor.pm-ai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent — runs the PM/AI profile once per day.

  Fires immediately on login/wake (RunAtLoad) and every 30 minutes thereafter.
  check_once.py is idempotent via the processed_videos JSON — videos already
  summarised are skipped instantly. The 30-min poll is a near-instant no-op on
  all but the first run of the day (no new videos = nothing to process).

  Previously ran only Mon/Wed/Fri — now runs daily so any day you power on
  you get the latest PM/AI content without waiting for a schedule slot.

  Install:
    cp launchd/com.youtube-monitor.pm-ai.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.youtube-monitor.pm-ai.plist

  Uninstall:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.youtube-monitor.pm-ai.plist

  Logs: /Users/makkarmandeep/stock-summary-youtube/cron.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube-monitor.pm-ai</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/makkarmandeep/stock-summary-youtube/run-pm-ai.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/makkarmandeep/stock-summary-youtube</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/makkarmandeep</string>
        <key>TZ</key>
        <string>America/Chicago</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- Fire on login/wake, then every 30 minutes. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/makkarmandeep/stock-summary-youtube/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/makkarmandeep/stock-summary-youtube/cron.log</string>
</dict>
</plist>
